# rancher-gce/nginx-dockprox/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
priority=1
command=service cron start
autorestart=true
startsecs=30

[program:nginx]
priority=10
command=nginx -g "daemon off;"
autorestart=true
startsecs=30
